# systemd/kanata-switcher.service
# Manual installation template for kanata-switcher systemd user service
# Copy to ~/.config/systemd/user/kanata-switcher.service and adjust ExecStart with correct port (if it's not 10000)

[Unit]
Description=Kanata layer switcher daemon
After=graphical-session.target
PartOf=graphical-session.target

[Service]
Type=simple
# Adjust the path and port as needed:
ExecStart=%h/.cargo/bin/kanata-switcher --quiet -p 10000
Restart=on-failure
RestartSec=5

[Install]
WantedBy=graphical-session.target
